# packaging/launchd/dev.bematist.collector.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Bematist collector — macOS LaunchAgent.
  Rendered at install time with @HOME@ and @BIN@ substituted for the user's
  $HOME and the resolved `bematist` binary path.

  Why LaunchAgent (not LaunchDaemon): the collector observes a developer's
  own sessions (Claude Code, Cursor, Codex live under $HOME). It must run
  as the user, triggered at GUI-login, not as root at boot.

  KeepAlive + RunAtLoad ensures the process restarts on crash and survives
  logout→login. SoftResourceLimits.Core=0 enforces CLAUDE.md §Security
  Rules (no crash dumps; matches the in-process harden() belt+suspenders).

  Environment: the `bematist` binary loads ~/.bematist/config.env itself
  via loadConfig() (apps/collector/src/config.ts). launchd just execs the
  binary directly. Earlier versions wrapped this in `/bin/sh -c` to source
  config.env first; that indirection was redundant AND broke on macOS
  machines where an EDR / Gatekeeper hook held the `sh` fork in SIGSTOP
  before exec (observed in the wild: ps STAT=T, 0-byte logs, no startup
  line). Keep this direct — the systemd counterpart
  (bematist.service.tmpl) uses ExecStart=@BIN@ serve for the same reason.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.bematist.collector</string>

  <key>ProgramArguments</key>
  <array>
    <string>@BIN@</string>
    <string>serve</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>@HOME@/.bematist/logs/out.log</string>

  <key>StandardErrorPath</key>
  <string>@HOME@/.bematist/logs/err.log</string>

  <key>WorkingDirectory</key>
  <string>@HOME@/.bematist</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>@HOME@</string>
  </dict>

  <key>SoftResourceLimits</key>
  <dict>
    <key>Core</key>
    <integer>0</integer>
  </dict>

  <key>HardResourceLimits</key>
  <dict>
    <key>Core</key>
    <integer>0</integer>
  </dict>

  <!-- Throttle restarts so a crash loop doesn't hog CPU. -->
  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
